# echo-unit@.service
[Unit]
Description=echo service

[Service]
ExecStart=-/bin/cat
StandardInput=socket

systemctl start echo.socket
systemctl stop echo.socket
